# p3-ai-noc.service
[Unit]
Description=P3 AI NOC TUI Dashboard on TTY2
After=network.target systemd-user-sessions.service getty@tty2.service
Conflicts=getty@tty2.service

[Service]
Type=simple
User=matty
Environment=TERM=xterm-256color
ExecStart=/usr/local/bin/p3ainoc
StandardInput=tty
StandardOutput=tty
StandardError=tty
TTYPath=/dev/tty2
TTYReset=yes
TTYVHangup=yes
Restart=always

[Install]
WantedBy=multi-user.target
